# crontab
* * * * * python3 /oceansuite/manage.py runscheduler 2>&1
